# chef/site-cookbooks/midas/templates/default/midas.upstart.erb
description "Start the midas node"

start on startup
stop on shutdown

respawn
respawn limit 10 5

script
su <%= @app_user %> -c "/bin/bash <<EOT
cd <%= @working_dir %>
node app.js --production
EOT"
end script
